# whm_bruteforce_monitoring.service
[Unit]
Description=WHM & cPanel Bruteforce Monitoring
After=network.target csf.service
Wants=csf.service

# Restart limits: max 5 restarts in 120 seconds before giving up
StartLimitIntervalSec=120
StartLimitBurst=5

[Service]
Type=simple
ExecStartPre=/usr/bin/test -f /etc/monitor/whm-bruteforce-monitor.sh
ExecStart=/bin/bash /etc/monitor/whm-bruteforce-monitor.sh
Restart=on-failure
RestartSec=10

# Resource limits
CPUQuota=20%
MemoryMax=128M

# Security hardening
ProtectHome=true
ProtectSystem=true
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=whm-bruteforce-monitor

[Install]
WantedBy=multi-user.target
